# bin/connect_network.service
[Unit]
Description=Check Network Exists or Setup Network
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/deadhead/deadstream/connect_network.py --debug 0
RemainAfterExit=yes
User=deadhead
[Install]
WantedBy=multi-user.target
