# bin/connect_network.service
[Unit]
Description=Check Network Exists or Setup Network
[Service]
Type=oneshot
ExecStart=/bin/sh -c 'while ! /usr/bin/python3 /home/deadhead/deadstream/connect_network.py --debug 0; do sleep 5; done'
RemainAfterExit=yes
User=deadhead
[Install]
WantedBy=multi-user.target
