# bin/connect_network.service
[Unit]
Description=Check Network Exists or Setup Network
After=multi-user.target
[Service]
Type=oneshot
RemainAfterExit=yes
User=steve
ExecStart=/usr/bin/python3 /home/steve/deadstream/connect_network.py --debug 0 
[Install]
WantedBy=multi-user.target
