# bin/connect_network.service
[Unit]
Description=Check Network Exists or Setup Network
[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/steve/deadstream/connect_network.py --debug 0 
RemainAfterExit=yes
User=steve
[Install]
WantedBy=multi-user.target
