# bin/serve_options.service
[Unit]
Description=Serve the Time Machine Options
Requires=connect_network.service
After=connect_network.service
[Service]
Type=simple
User=deadhead
ExecStart=/usr/bin/python3 /home/deadhead/deadstream/serve_options.py --debug 0
RemainAfterExit=no
Restart=on-failure
[Install]
WantedBy=multi-user.target
