# bin/timemachine.service
[Unit]
Description=Grateful Dead Time Machine
#Wants=connect_network.service
Requires=connect_network.service
After=connect_network.service
#BindsTo=check_network.service
[Service]
Type=simple
User=steve
ExecStart=/usr/bin/python3 /home/steve/deadstream/main.py --debug 0 --box v1
#Restart=on-failure
[Install]
WantedBy=multi-user.target
