# bin/timemachine.service
[Unit]
Description=Grateful Dead Time Machine
#Wants=connect_network.service
Requires=connect_network.service
After=connect_network.service
#BindsTo=check_network.service
[Service]
Type=simple
User=deadhead
ExecStart=/usr/bin/python3 /home/deadhead/deadstream/main.py
#Restart=on-failure
[Install]
WantedBy=multi-user.target
