# timemachine/bin/calibrate.service
[Unit]
Description=Calibrate Time Machine
[Service]
Type=oneshot
User=deadhead
ExecStart=/bin/bash -c "source /home/deadhead/timemachine/bin/activate && calibrate"
ExecStart=/bin/bash -c "source /home/deadhead/timemachine/bin/activate && stand_by.sh"
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
