# timemachine/bin/connect_network.service
[Unit]
Description=Check Network Exists or Setup Network
Wants=network.target
After=systemd-user-sessions.service
[Service]
Type=oneshot
User=deadhead
#ExecStart=/bin/sh -c 'while ! /usr/bin/python3 /home/deadhead/deadstream/timemachine/connect_network.py --debug 0; do sleep 5; done'
ExecStart=/home/deadhead/timemachine/bin/connect_network --debug 0
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
