# timemachine/bin/update.service
[Unit]
Description=Update the Grateful Dead Time Machine
Requires=connect_network.service
[Service]
Type=simple
User=deadhead
ExecStart=/home/deadhead/deadstream/bin/update.sh
[Install]
WantedBy=multi-user.target
